# reversi/effective_reversi.ini
[program:effective_reversi]
directory=/home/ec2-user/effective_reversi/reversi
command=/usr/bin/pipenv run daphne reversi.routing:application
autostart=true
autorestart=true
user=ec2-user
environment=HOME="/home/ec2-user",LD_LIBRARY_PATH="/usr/local/lib"
stdout_logfile=/var/log/supervisor/effective_reversi.log
stdout_logfile_maxbites=1MB
redirect_stderr=true
